# scripts/selenium-hub.conf
author "ravi"
description "selenium grid hub copy to ~/.config/upstart"

start on runlevel [2345]

respawn

chdir /srv/moz-grid-config/

exec ant launch-hub

